# crontab/everyday.sh
0 12 * * * cd /var/opt/comic/dev/crawler && source activate manga && scrapy crawl natalie >> natalie_crawl.log 2>&1 >/dev/null 2>&1